# projects/Clearwater/Docker/clearwater-docker/astaire/crontab
*/10 * * * * root /bin/bash -c "/usr/share/clearwater/bin/scale_astaire.sh"
# the above given line will run scale down script every 2 minutes
