# Clearwater/Docker/clearwater-docker/base/supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command=bash -c "monit & monit restart all"


